# docker/node_app/supervisor.conf
[supervisord]
nodaemon=true

[program:app1]
command=node /opt/techtalk-2014-dec/app/server.js
environment=PORT="3000"
stdout_logfile=/logs/app.log
stdout_logfile_maxbytes=1GB
redirect_stderr=true

